# supervisor/conf.d/crontab.ini
[program:crontab]
command = /usr/sbin/crond -f
autostart=true
autorestart=true
priority=5

stdout_events_enabled=true
stdout_logfile_maxbytes=51200
stdout_logfile=/var/log/crontab-stdout.log

stderr_events_enabled=true
stderr_logfile_maxbytes=51200
stderr_logfile=/var/log/crontab-stderr.log


