# supervisor/conf.d/openresty.ini
[program:openresty]
command=/usr/local/openresty/bin/openresty -g "daemon off; error_log /dev/stderr info;"
autostart=true
autorestart=true
priority=10

stdout_events_enabled=true
stdout_logfile_maxbytes=51200
stdout_logfile=/var/log/openresty-stdout.log

stderr_events_enabled=true
stderr_logfile_maxbytes=51200
stderr_logfile=/var/log/openresty-stderr.log
