# crontab/test.crontab
*/1 * * * * echo `date` >> /tmp/test.log